# deps/sv/void-nfs-repos/finish
#!/bin/sh

[ -r ./conf ] && . ./conf

ln -sf /usr/share/xbps.d /etc/xbps.d
umount ${LOCAL_MOUNT:=/mnt/repos}
