# deps/sv/void-nfs-repos/finish
#!/bin/sh

[ -r ./conf ] && . ./conf

rm -rf /etc/xbps.d
umount ${LOCAL_MOUNT:=/mnt/repos}
