# deps/sv/void-nfs-repos/run
#!/bin/sh

[ -r ./conf ] && . ./conf

[ ! -d ${LOCAL_MOUNT:=/mnt/void-nfs-repos} ] && mkdir $LOCAL_MOUNT

if [ -d ${NFS_CONF_DIR:=/usr/local/share/xbps.nfs.d} ]; then
  ln -sf $NFS_CONF_DIR /etc/xbps.d
else
  mkdir -p $NFS_CONF_DIR

  repository="repository=${LOCAL_MOUNT}/current"
  [ $MUSL = 'true' ] && repository="${repository}/musl"
  echo "$repository" > $NFS_CONF_DIR/00-repository-main.conf

  [ $NONFREE = 'true' ] && nonfree="${repository}/nonfree"
  echo "$nonfree" > $NFS_CONF_DIR/10-repository-nonfree.conf
fi

if [ $NONFREE != 'true' ] && [ -f $NFS_CONF_DIR/10-repository-nonfree.conf ]; then
  rm -rf $NFS_CONF_DIR/10-repository-nonfree.conf
fi

ln -sf $NFS_CONF_DIR /etc/xbps.d

mount -t nfs ${NFS_REPO_IP:=10.20.30.2}:${NFS_REPO_DIR:=/srv/void} $LOCAL_MOUNT
exec pause
